# launchd/com.yannickherrero.opencode-web.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.yannickherrero.opencode-web</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/yannickherrero/.opencode/bin/opencode</string>
    <string>web</string>
    <string>--port</string>
    <string>4096</string>
    <string>--hostname</string>
    <string>127.0.0.1</string>
    <string>--cors</string>
    <string>https://macbook.tailf3d9b7.ts.net</string>
    <string>--print-logs</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/yannickherrero/dev</string>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>5</integer>

  <key>StandardOutPath</key>
  <string>/Users/yannickherrero/Library/Logs/opencode-web.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/yannickherrero/Library/Logs/opencode-web.log</string>
</dict>
</plist>
