# launchd/com.yannickherrero.opencode-web-menu.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.yannickherrero.opencode-web-menu</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/yannickherrero/Applications/OpenCode Status.app/Contents/MacOS/OpenCodeStatus</string>
  </array>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
    <key>Crashed</key>
    <true/>
  </dict>
  <key>ThrottleInterval</key>
  <integer>5</integer>

  <key>StandardOutPath</key>
  <string>/Users/yannickherrero/Library/Logs/opencode-web-menu.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/yannickherrero/Library/Logs/opencode-web-menu.log</string>
</dict>
</plist>
